# NON-WINDOWS NON-C/LINUX/jobs.cron
0,5,10,15,20,25,30,35,40,45,50,55 * * * * cd /home/ec2-user/ChatScript/BINARIES && ./LinuxChatScript64 >> cronserver.log 2>&1
